# scheduler/crontab
0 2 * * * /bin/bash -c 'sleep $((RANDOM % 3600)) && curl -X POST http://api:8000/run/crawler'
